# bot_rasp_backup.timer
[Unit]
Description=Daily Telegram Bot Raspisanie backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true

[Install]
WantedBy=timers.target
