# bot_rasp_web.service
[Unit]
Description=Telegram Bot Raspisanie (web panel)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bot_rasp
EnvironmentFile=/home/bot_rasp/.env
ExecStart=/home/bot_rasp/venv/bin/python api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
