# bot_rasp_web.service
[Unit]
Description=Telegram Bot Raspisanie (web panel)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bot_rasp
EnvironmentFile=/home/bot_rasp/.env
ExecStart=/home/bot_rasp/venv/bin/gunicorn -w 2 -b 127.0.0.1:5001 api:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
