# bot_rasp_cleanup.timer
[Unit]
Description=Weekly Telegram Bot Raspisanie technical table cleanup

[Timer]
OnCalendar=Sun 04:10:00
Persistent=true

[Install]
WantedBy=timers.target
